# fix-share.conf
start on startup
task
exec /usr/bin/fix-share
